# overlay/etc/init.d/brick
#!/sbin/openrc-run

description="Brick iot app"

required_dirs="/etc/brick"

command="/media/mmcblk0p1/opt/usr/bin/brick"
command_args="--config-dir /etc/brick"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"
output_log="/var/log/brick"
error_log="/var/log/brick"
